# devbox-dockerfile/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:vncserver]
priority=10
directory=/root
command=/usr/bin/vncserver :1 -geometry 1280x800 -depth 24 -localhost no -xstartup /root/.vnc/xstartup
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/vncserver.stderr.log
stdout_logfile=/var/log/supervisor/vncserver.stdout.log
stopasgroup=true
killasgroup=true

[program:novnc]
priority=20
directory=/usr/share/novnc
command=/usr/bin/websockify --web=/usr/share/novnc 8080 localhost:5901
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/novnc.stderr.log
stdout_logfile=/var/log/supervisor/novnc.stdout.log

[program:code-server]
priority=30
directory=/root
command=/usr/bin/code-server --bind-addr 0.0.0.0:8443 --auth none
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/code-server.stderr.log
stdout_logfile=/var/log/supervisor/code-server.stdout.log

[program:redis]
priority=40
command=/usr/bin/redis-server
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/redis.stderr.log
stdout_logfile=/var/log/supervisor/redis.stdout.log

[program:mariadb]
priority=50
command=/usr/bin/mysqld_safe
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/mariadb.stderr.log
stdout_logfile=/var/log/supervisor/mariadb.stdout.log[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:vncserver]
priority=10
directory=/root
command=/usr/bin/vncserver :1 -geometry 1280x800 -depth 24 -localhost no -xstartup /root/.vnc/xstartup
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/vncserver.stderr.log
stdout_logfile=/var/log/supervisor/vncserver.stdout.log
stopasgroup=true
killasgroup=true

[program:novnc]
priority=20
directory=/usr/share/novnc
command=/usr/bin/websockify --web=/usr/share/novnc 8080 localhost:5901
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/novnc.stderr.log
stdout_logfile=/var/log/supervisor/novnc.stdout.log

[program:code-server]
priority=30
directory=/root
command=/usr/bin/code-server --bind-addr 0.0.0.0:8443 --auth none
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/code-server.stderr.log
stdout_logfile=/var/log/supervisor/code-server.stdout.log

[program:redis]
priority=40
command=/usr/bin/redis-server
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/redis.stderr.log
stdout_logfile=/var/log/supervisor/redis.stdout.log

[program:mariadb]
priority=50
command=/usr/bin/mysqld_safe
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/mariadb.stderr.log
stdout_logfile=/var/log/supervisor/mariadb.stdout.log[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:vncserver]
priority=10
directory=/root
command=/usr/bin/vncserver :1 -geometry 1280x800 -depth 24 -localhost no -xstartup /root/.vnc/xstartup
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/vncserver.stderr.log
stdout_logfile=/var/log/supervisor/vncserver.stdout.log
stopasgroup=true
killasgroup=true

[program:novnc]
priority=20
directory=/usr/share/novnc
command=/usr/bin/websockify --web=/usr/share/novnc 8080 localhost:5901
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/novnc.stderr.log
stdout_logfile=/var/log/supervisor/novnc.stdout.log

[program:code-server]
priority=30
directory=/root
command=/usr/bin/code-server --bind-addr 0.0.0.0:8443 --auth none
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/code-server.stderr.log
stdout_logfile=/var/log/supervisor/code-server.stdout.log

[program:redis]
priority=40
command=/usr/bin/redis-server
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/redis.stderr.log
stdout_logfile=/var/log/supervisor/redis.stdout.log

[program:mariadb]
priority=50
command=/usr/bin/mysqld_safe
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/mariadb.stderr.log
stdout_logfile=/var/log/supervisor/mariadb.stdout.log